# desktop/config-files/etc/NetworkManager/dispatcher.d/10-update-timesyncd
#!/usr/bin/env dash 

set -e

INTERFACE=$1
ACTION=$2
_conf_dir=/etc/systemd/timesyncd.conf.d
_timesyncd_state="$( systemctl show --property "ActiveState" systemd-timesyncd.service )"

# Do we have a connection ID?
if [ -z "$CONNECTION_UUID" ]; then
    exit 0
fi

# Is the systemd-timesyncd.service active?
if [ "$_timesyncd_state" != "ActiveState=active" ]
then
    # timesyncd is not active.
    exit 0
fi

if [ ! -d "$_conf_dir" ]; then

    # Creating timesyncd overlay configuration directory
    mkdir -p "$_conf_dir"
fi

case $ACTION in
    up | vpn-up | dhcp4-change | dhcp6-change)

        # Network just went up or configuration received by DHCP has changed.
    
        # Did we receive any NTP server addresses by DHCP?
        if [ -z "$DHCP4_NTP_SERVERS" ] && [ -z "$DHCP6_NTP_SERVERS" ]; then

            # Neither IPv6 nor IPv4 NTP server address received. 
            exit 0
        fi

        if [ -n "$DHCP4_NTP_SERVERS" ]; then
            _ntp_servers="$_ntp_servers $DHCP4_NTP_SERVERS"
        fi
        if [ -n "$DHCP6_NTP_SERVERS" ]; then
            _ntp_servers="$_ntp_servers $DHCP6_NTP_SERVERS"
        fi

        exec > "${_conf_dir}/${CONNECTION_UUID}.conf"
        echo "# Generated by $(basename "$0")for $INTERFACE"
        echo "[Time]"
        echo "NTP=${_ntp_servers}"
        systemctl restart systemd-timesyncd
        ;;
    down | vpn-down)

        # Connection just went down, removing its timesyncd configuration."
        rm -f "${_conf_dir}/${CONNECTION_UUID}.conf"
        systemctl restart systemd-timesyncd
        ;;
esac
